# deploy/nsdiggest.timer
[Unit]
Description=NSDiggest timer (disabled)

[Timer]
# Intentionally cleared to disable scheduled runs.
OnCalendar=
Persistent=true
RandomizedDelaySec=5min
Unit=nsdiggest.service

[Install]
WantedBy=timers.target
