# deploy/nsdiggest.timer
[Unit]
Description=NSDiggest daily timer (07:00 Europe/Warsaw)

[Timer]
OnCalendar=*-*-* 07:00:00 Europe/Warsaw
Persistent=true
RandomizedDelaySec=5min
Unit=nsdiggest.service

[Install]
WantedBy=timers.target
